# packaging/ubuntu/docker.upstart
description     "Run docker"

start on started network-services
stop on stopping network-services

respawn

exec /usr/bin/docker -d
